# buildroot_cfg/overlay_fs/etc/inittab
::sysinit:/etc/init.d/xorg.rc
::restart:/sbin/init